# deploy/micromania-deals.service
[Unit]
Description=Micromania deals watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/micromania-deals
EnvironmentFile=/etc/micromania-deals.env
ExecStart=/usr/bin/python3 /opt/micromania-deals/micromania_deals.py
Restart=always
RestartSec=10
# Redémarre proprement si le process se fige
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
