# deploy/micromania-deals.service
[Unit]
Description=Micromania deals watcher
After=network-online.target
Wants=network-online.target

[Service]
# notify => le script envoie un heartbeat (sd_notify) ; si systemd ne le reçoit
# plus pendant WatchdogSec, il considère le process "figé" et le redémarre.
Type=notify
WatchdogSec=600
WorkingDirectory=/opt/micromania-deals
EnvironmentFile=/etc/micromania-deals.env
ExecStart=/opt/micromania-deals/venv/bin/python /opt/micromania-deals/micromania_deals.py
Restart=always
RestartSec=10
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
